# ops/backup/careflow-restore-drill.timer
[Unit]
Description=careflow-lite: run the restore drill every week

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
