# ops/backup/careflow-backup.service
[Unit]
Description=careflow-lite: nightly MongoDB backup (whole instance)
After=network-online.target mongod.service
Wants=network-online.target

[Service]
Type=oneshot
User=roger
ExecStart=/home/roger/apps/careflow-lite/ops/backup/mongo-backup.sh
# Locking lives inside the script (it re-execs itself under flock), so a manual run and the timer
# cannot collide either. An ExecStartPre flock would release the lock before ExecStart even began.
